# cron.d/imaging-rename
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/opt/Catmandu/bin
MAILTO=root
HOME=/
PERL5LIB=/opt/Catmandu/lib:/opt/Catmandu-MARC/lib:/opt/Imaging/lib
*/4 2-22 * * * root /opt/Imaging/perl/cron-rename.pl 1>> /var/log/Imaging/imaging-rename.log 2>> /var/log/Imaging/imaging-rename.error
