# cron.d/imaging-backup
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
MAILTO=imaging
HOME=/home/imaging

0 0 * * * imaging /opt/Imaging/bin/backup.pl --dir /mnt/data01/backup --database imaging 1>> /var/log/Imaging/backup.log 2>> /var/log/Imaging/backup.error
