# cron.d/imaging-status
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/opt/Catmandu/bin
MAILTO=imaging
HOME=/
PERL5LIB=/opt/Catmandu/lib:/opt/Catmandu-MARC/lib:/opt/Imaging/lib
*/4 02-21 * * * imaging /opt/Imaging/bin/cron-status.pl 1>> /var/log/Imaging/cron-status.log 2>> /var/log/Imaging/cron-status.error
