# cron.d/imaging-register
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/opt/Catmandu/bin
MAILTO=root
HOME=/
PERL5LIB=/opt/Catmandu/lib:/opt/Catmandu-MARC/lib:/opt/Imaging/lib
30 22 * * * root /opt/Imaging/perl/delete_scans_done.pl && /opt/Imaging/perl/cron-register.pl 1>> /var/log/Imaging/cron-register.log 2>> /var/log/Imaging/cron-register.error
