# cron.d/imaging-nightly
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
MAILTO=imaging
HOME=/

IMAGING_HOME=/opt/imaging
IMAGING_LOG=/var/log/imaging

10 19 * * * imaging /opt/imaging/bin/imaging-nightly
