# cron.d/imaging-check
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/opt/Catmandu/bin
MAILTO=root
HOME=/
PERL5LIB=/opt/Catmandu/lib:/opt/Catmandu-MARC/lib:/opt/Imaging/lib
*/2 02-22 * * * root /opt/Imaging/perl/cron-check.pl 1>> /var/log/Imaging/cron-check.log 2>> /var/log/Imaging/cron-check.error
